# debian/dsched.dsched.service
[Unit]
Description=process scheduler manager
After=phosh.service

[Service]
Type=simple
ExecStart=/usr/bin/dsched
Restart=always
RestartSec=5s
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
